# preentrega_03/supervisord.conf
[supervisord]
nodaemon=true

[program:airflow-webserver]
command=airflow webserver --port 8080
autostart=true
autorestart=true
stdout_logfile=/usr/local/airflow/logs/webserver.log
stderr_logfile=/usr/local/airflow/logs/webserver.err

[program:airflow-scheduler]
command=airflow scheduler
autostart=true
autorestart=true
stdout_logfile=/usr/local/airflow/logs/scheduler.log
stderr_logfile=/usr/local/airflow/logs/scheduler.err